# .config/systemd/user/batsignal.service
[Unit]
Description=Lightweight battery daemon

[Service]
ExecStart=/usr/bin/batsignal -w 18 -c 8 -d 3

[Install]
WantedBy=default.target
